# apps/alert-agent/manifests/files/.crontab
# supercronic reads /home/agent/.crontab (baked s6 service). Env (PYTHONPATH,
# secrets) is inherited from the agent container. Surge gate every 15m; digest daily 08:00.
*/15 * * * * /opt/alert-agent-bin/surge-gate
0 8 * * * /opt/alert-agent-bin/digest
